# Ejecutables/inicioadmplus.service
[Unit]
Description=ADMPlus

[Service]
Type=forking
ExecStart=/bin/autoinicios
StandardOutput=null

[Install]
WantedBy=multi-user.target
